# deploy/tesla-journey-os.service
[Unit]
Description=Tesla Journey OS — Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/tesla-journey-os
Environment=PYTHONPATH=/opt/tesla-journey-os/backend
ExecStart=/opt/tesla-journey-os/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/tesla-journey-os/data
ReadOnlyPaths=/opt/tesla-journey-os/backend /opt/tesla-journey-os/config.yaml

[Install]
WantedBy=multi-user.target
